# server/deploy/chatgpt.service
# /etc/systemd/system/chatgpt.service
[Unit]
Description=ChatGPT on Aliyun (Node API)
After=network.target

[Service]
Type=simple
Environment=PORT=3000
Environment=OPENAI_MODEL=gpt-4o-mini
Environment=CORS_ORIGIN=https://YOUR_DOMAIN
WorkingDirectory=/opt/chatgpt-on-aliyun
ExecStart=/usr/bin/node server/index.js
Restart=always
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
